# systemd/voxtype-pill.service
[Unit]
Description=VoxType pill overlay

[Service]
Type=simple
ExecStart=%h/.local/bin/voxtype-pill
Restart=on-failure
RestartSec=3
