# systemd/voxtyped.service
[Unit]
Description=VoxType daemon (hotkey detection for voice typing)
Wants=voxtype-ydotoold.service voxtype-pill.service
After=voxtype-ydotoold.service

[Service]
Type=simple
ExecStart=%h/.local/bin/voxtyped
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
